# ops/systemd/alphatrend-demand-refresh.service
[Unit]
Description=AlphaTrend V3 economic demand refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/alphatrend
EnvironmentFile=/srv/alphatrend/.env
ExecStart=/usr/bin/env npm run job:demand -- --all
